# assets/ALL/ALL/home/aki/.config/autostart/vesktop.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Vesktop
Comment=Vesktop autostart script
Exec=/usr/bin/vesktop
Terminal=false
StartupNotify=false